# .config/systemd/user/brightness-manager.service
[Unit]
Description=Auto Brightness Manager
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/brightness-manager
Nice=10
